# infra/supervisor/supervisord.conf
; =============================================================================
; OpenPush — Supervisor config for the `queue` service
; -----------------------------------------------------------------------------
; Runs a pool of `queue:work push` workers plus a horizon-free failsafe. Each
; `queue` container runs WORKERS processes; scale further with
; `docker compose up --scale queue=N` (total workers = N * numprocs).
; =============================================================================

[supervisord]
nodaemon=true                 ; run in foreground as PID 1 of the container
; The container already runs as the unprivileged `openpush` user (see the
; backend Dockerfile), so supervisord itself runs unprivileged and does NOT
; switch users — hence no `user=` here and none on the programs below.
logfile=/dev/stdout           ; ship logs to the container stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid  ; /tmp is writable by the non-root user

; -----------------------------------------------------------------------------
; push queue workers
; -----------------------------------------------------------------------------
; API_CONTRACT lifecycle: SendFcmBatchJob lands on the `push` queue. These
; workers drain it. Tuning:
;   --queue=push        only process the push queue (dedicated throughput)
;   --sleep=1           poll delay when idle
;   --tries=3           matches the contract's max-3-retries policy
;   --backoff=10,30,60  exponential-ish backoff between attempts
;   --max-time=3600     recycle the worker hourly to release memory
;   --max-jobs=1000     recycle after N jobs (leak guard)
; stopwaitsecs must exceed the longest job so SIGTERM lets a batch finish
; (graceful shutdown; requires the pcntl extension baked into the image).
[program:push-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=push --sleep=1 --tries=3 --backoff=10,30,60 --max-time=3600 --max-jobs=1000 --no-interaction
directory=/var/www/html
numprocs=%(ENV_QUEUE_WORKERS)s ; per-container worker count (default set below via env)
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=120               ; allow an in-flight FCM batch to complete
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; -----------------------------------------------------------------------------
; default queue worker (transactional: mail, audit logs, misc jobs)
; -----------------------------------------------------------------------------
[program:default-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600 --no-interaction
directory=/var/www/html
numprocs=1
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=60
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
